# monitor_test.service
[Unit]
Description=Monitor test process
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash /usr/local/bin/monitor_test.sh

[Install]
WantedBy=multi-user.target